# elk/elk_kafka_service_makager_script/kafka/kafkad
#/bin/bash
# chkconfig: 2345 20 81
# description: kafka service manager script
# auth: guomaoqiu
# date: 2015-10-09
CONFIG=/usr/local/kafka

#
if [ ! -d  ${CONFIG}/nohup/ ];then
    mkdir  ${CONFIG}/nohup/{zookeeper,kafka} -p
fi

#
start(){
    netstat -nlpt | grep "9092" >/dev/null
    if [[ $? -eq 0 ]]; then
      echo "Kafka already start,nothing to do"
      exit
    fi
    cd ${CONFIG}/nohup/zookeeper
    echo  "Starting kafka Service"
    nohup /bin/bash /usr/local/kafka/bin/kafka-server-start.sh /usr/local/kafka/config/server.properties 2>&1 &
}

stop(){
    echo "Stoping kafka"
    ps ax | grep -i 'server.properties' | grep -v grep | awk '{print $1}' | xargs kill -SIGTERM
}

status(){
    netstat -nlpt | grep "9092" >/dev/null
    if [ $? -eq 0 ];then
       echo "Kafka is running"
    else
       echo "Kafka is stopped" 
    fi
}


case "$1" in
start)
  start
  ;;
stop)
  stop
  ;;
status)
  status
  ;;
esac
